# src/etc/supervisor.d/gate.ini
[supervisord]
nodaemon=true

[program:redsocks]
command=/usr/bin/redsocks -c /etc/redsocks.conf
autorestart=true

[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon --conf-dir=/conf/dnsmasq.d --conf-file=/etc/dnsmasq.conf
autorestart=true

[program:pdnsd]
command=/usr/sbin/pdnsd -c /etc/pdnsd.conf
autorestart=true

[program:chinadns]
command=/usr/local/bin/chinadns -c /usr/local/share/chnroute.txt -s 202.101.172.35,127.0.0.1:5302 -p 5301 -v
autorestart=true

[program:ssr]
command=/usr/bin/python /shadowsocks/shadowsocks/local.py -c /conf/ssr.conf
autorestart=true

[program:cron]
command = crond -f
autostart=true
autorestart=true
